# sysadmin/heltour.conf
#!upstart
description "heltour gunicorn server"
author      "Lakin Wecker"

start on (started networking)
stop on shutdown

script
    export HOME="/var/www/heltour.lakin.ca/"

    exec sudo -u www-data /var/www/heltour.lakin.ca/run-heltour.sh 8580
end script
